# supervisord-debug.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisor/supervisord.pid

[program:fair-projects]
command=authbind python3 run.py --debug --ip 0.0.0.0 --port 80
user=www-data
directory=%(ENV_APPDIR)s
logfile=/var/log/supervisor/%(ENV_SERVERNAME)s-supervisord.log
pidfile=/var/run/supervisor/%(ENV_SERVERNAME)s-supervisord.pid
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0